# rc.local
#!/bin/sh -e
#
# rc.local

sudo service cron start &

echo "\e[92m
--> Starting server\e[0m"
iptables -t nat -A PREROUTING -p tcp --dport 80 -j REDIRECT --to-port 17325
su USERR -c 'cd HOMEE/git/fuettr_prototype/Web_Application/Webserver/ng2 && npm start &'

_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  echo "
\e[41m My IP address is \e[30m$_IP\e[0m\n"
fi


exit 0
