# rc.local
#!/bin/sh -e
#
# rc.local

sudo service cron start &

echo "
--> Starting server"
su USER -c 'cd USERR/git/fuettr_prototype/Web_Application/Webserver/ng2 && npm start &'

_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  echo -e "My IP address is \e[92m$_IP\e[0m"
fi


exit 0
