# rc.local
#!/bin/sh -e
#
# rc.local

sudo service cron start &

echo "\e[92m
--> Starting server\e[0m"
su USER -c 'cd USERR/git/fuettr_prototype/Web_Application/Webserver/ng2 && npm start &'

_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  echo "\e[41mMy IP address is \e[30m$_IP\e[0m"
fi


exit 0
